# nyasakura.service
[Unit]
Description=NyaSakura Frontend Server
After=network.target

[Service]
ExecStart=/usr/bin/node /path/to/nyasakura/nyasakura-service.js
# 确保node进程接收到正确的停止信号
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30
# 确保服务停止时所有相关进程都被终止
SendSIGKILL=yes
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=nyasakura-frontend
User=your-username
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 